# files/system/usr/lib/systemd/system/nix-maintenance.timer
[Unit]
Description=Weekly Nix maintenance

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
